# conf/crontab
SHELL=/bin/sh
*/15 * * * *       /home/patches/update-patches.sh
